# configFiles/alarm/alarm.timer
[Unit]
Description=Alarm

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 06:35

[Install]
WantedBy=multi-user.target
