# configFiles/apache/systemd-ampacheupdate.timer
[Unit]
Description=Update songs from YouTube

[Timer]
OnCalendar=*-*-* 05:50:00

[Install]
WantedBy=multi-user.target
